# gunicorn/gunicorn_dksharma.socket
[Unit]
Description=gunicorn socket for dksharma

[Socket]
ListenStream=/run/gunicorn_dksharma.sock

[Install]
WantedBy=sockets.target